# docker/root/etc/s6-overlay/s6-rc.d/bindfs/run
#!/command/execlineb -P
redirfd -w 1 /dev/null
redirfd -w 2 /dev/null
bindfs -f -o allow_other --resolve-symlinks --map=root/www-data /mounts /Video_Library
