# docker/root/etc/s6-overlay/s6-rc.d/bindfs/run
#!/command/execlineb -P
sh -c "bindfs -f -o allow_other --resolve-symlinks --map=root/www-data /mounts /Video_Library > /jellygrail/log/bindfs 2>&1"
