# docker/root/etc/s6-overlay/s6-rc.d/kodi_bindfs/run
#!/command/execlineb -P
sh -c "bindfs -f -o allow_other --resolve-symlinks --display-nfo --map=root/www-data /mounts /Kodi_Video_Library > /jellygrail/log/kodi_bindfs 2>&1"
