# docker/root/etc/s6-overlay/s6-rc.d/jellyfin/run
#!/command/execlineb -P
redirfd -w 1 /dev/null
redirfd -w 2 /dev/null
s6-setuidgid www-data
export JELLYFIN_FFMPEG__NOVALIDATION true
export JELLYFIN_FFmpeg__analyzeduration 6000000
export LD_PRELOAD /usr/lib/x86_64-linux-gnu/libjemalloc.so.2
/jellyfin/jellyfin --datadir /jellygrail/jellyfin/config --cachedir /jellygrail/jellyfin/cache --ffmpeg /usr/lib/jellyfin-ffmpeg/ffmpeg
