# docker/root/etc/s6-overlay/s6-rc.d/jellyfin/run
#!/command/execlineb -P
redirfd -w 1 /dev/null
redirfd -w 2 /dev/null
s6-setuidgid www-data
env JELLYFIN_FFmpeg__analyzeduration=4000000 /jellyfin/jellyfin --datadir /jellygrail/jellyfin/config --cachedir /jellygrail/jellyfin/cache --ffmpeg /usr/lib/jellyfin-ffmpeg/ffmpeg
