# docker/root/etc/s6-overlay/s6-rc.d/mariadb/run
#!/usr/bin/execlineb -P
s6-setuidgid mysql
foreground { /usr/bin/mysql_install_db --user=mysql --datadir=/var/lib/mysql }
foreground { /usr/bin/mysqld_safe --datadir=/var/lib/mysql }
