# docker/root/etc/s6-overlay/s6-rc.d/kodi_addons/up
#!/command/execlineb -P
/bin/bash /kodi_addons.sh