# docker/root/etc/s6-overlay/s6-rc.d/init-env/up
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
export MYSQL_ROOT_PASSWORD=lgclvdnmj1492
export MYSQL_DATABASE=kodi_video131
export MYSQL_USER=kodi
export MYSQL_PASSWORD=kodi
export DATADIR=/jellygrail/data/mysql
