# docker/root/etc/s6-overlay/s6-rc.d/mariadbd/run
#!/command/execlineb -P
s6-setuidgid mysql
redirfd -w 1 /dev/null
redirfd -w 2 /dev/null
mariadbd
