# docker/root/etc/s6-overlay/s6-rc.d/cache_check_bindfs/run
#!/command/execlineb -P
bindfs -f -o allow_other --resolve-symlinks --map=root/www-data /mounts /Cache_Check_Video_Library
