# docker/root/etc/s6-overlay/s6-rc.d/jelly_update/run
#!/command/execlineb -P
python3 /usr/local/bin/grail_jelly/main.py