# docker/root/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
redirfd -w 1 /dev/null
redirfd -w 2 /dev/null
nginx -g "daemon off;"