# docker/root/etc/s6-overlay/s6-rc.d/mariadb_install/up
#!/command/execlineb -P
/bin/bash /install_mariadb.sh
