# docker/root/etc/s6-overlay/s6-rc.d/init_mariadb/up
#!/command/execlineb -P
/bin/bash /init_mariadb.sh
